# dist/etc/udev/rules.d/10-disable-ehci.rules
ACTION=="add", SUBSYSTEM=="pci", DRIVER=="ehci_hcd", RUN+="/bin/sh -c 'echo -n %k > %S%p/driver/unbind'"
