# dist/etc/skel/.config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=xbindkeys
Exec=xbindkeys
StartupNotify=false
NoDisplay=true
X-GNOME-Autostart-Phase=Initialization
X-KDE-autostart-phase=1
